# poky/meta/recipes-core/initscripts/initscripts-1.0/wifi.service
[Unit]
Description=start initscript upon first boot

[Service]
Type=simple
ExecStart=/bin/sh -c 'sleep 5 ; etc/init.d/setup-wifi.sh'
