# meta-train/recipes-systemd/systemd/files/wifi.service
[Unit]
Description=Wifi setup

[Service]
ExecStart=/etc/init.d/setup-wifi.sh

[Install]
WantedBy=multi-user.target
